# packages/podman-chrome-streaming/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
loglevel=info

; ── 1. Runtime (repo clone + funny server) ─────────────────────────────────────
; This is the main process — clones the repo, then starts the published Funny server runtime.
; It must complete repo setup before the server becomes available.
[program:runtime]
command=/root/.bun/bin/bun run /app/packages/podman-chrome-streaming/src/runtime/entrypoint.ts
environment=
    REPO_MODE="%(ENV_REPO_MODE)s",
    REPO_URL="%(ENV_REPO_URL)s",
    REPO_REF="%(ENV_REPO_REF)s",
    WORK_BRANCH="%(ENV_WORK_BRANCH)s",
    GIT_TOKEN="%(ENV_GIT_TOKEN)s",
    GIT_TOKEN_FILE="%(ENV_GIT_TOKEN_FILE)s",
    GIT_USERNAME="%(ENV_GIT_USERNAME)s",
    WORKSPACE_PATH="%(ENV_WORKSPACE_PATH)s",
    FUNNY_PORT="%(ENV_FUNNY_PORT)s",
    CLIENT_ORIGIN="%(ENV_CLIENT_ORIGIN)s",
    AUTH_MODE="%(ENV_AUTH_MODE)s",
    FUNNY_DATA_DIR="%(ENV_FUNNY_DATA_DIR)s",
    FUNNY_SERVER_BIN="%(ENV_FUNNY_SERVER_BIN)s",
    ENABLE_RUNTIME="%(ENV_ENABLE_RUNTIME)s",
    ENABLE_STREAMING="%(ENV_ENABLE_STREAMING)s"
autostart=%(ENV_ENABLE_RUNTIME)s
priority=5
autorestart=true
startsecs=8
startretries=10
stopwaitsecs=20
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/runtime.log
stderr_logfile=/var/log/supervisor/runtime.log

; ── 2. Virtual display ─────────────────────────────────────────────────────────
[program:xvfb]
command=Xvfb :99 -screen 0 1280x720x24 -ac -nolisten tcp
priority=10
autorestart=true
autostart=%(ENV_ENABLE_STREAMING)s
startretries=5
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.log

; ── 3. Chrome (non-headless on virtual display for noVNC view) ─────────────────
[program:chrome]
environment=DISPLAY=":99"
command=chromium
    --no-sandbox
    --disable-dev-shm-usage
    --disable-gpu
    --window-size=1280,720
    --remote-debugging-port=%(ENV_CHROME_DEBUG_PORT)s
    --remote-debugging-address=0.0.0.0
    --no-first-run
    --no-default-browser-check
    --disable-extensions
    --disable-background-networking
    --disable-sync
    --disable-translate
    about:blank
priority=20
autorestart=true
autostart=%(ENV_ENABLE_STREAMING)s
startsecs=4
startretries=10
stopwaitsecs=15
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/chrome.log
stderr_logfile=/var/log/supervisor/chrome.log

; ── 4. x11vnc (exposes the virtual display as VNC) ────────────────────────────
[program:x11vnc]
command=x11vnc -display :99 -nopw -listen 0.0.0.0 -xkb -forever -shared
priority=30
autorestart=true
autostart=%(ENV_ENABLE_STREAMING)s
startsecs=3
startretries=10
stopwaitsecs=10
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.log

; ── 5. noVNC (browser-accessible VNC) ─────────────────────────────────────────
[program:novnc]
command=websockify --web=/usr/share/novnc %(ENV_NOVNC_PORT)s localhost:5900
priority=40
autorestart=true
autostart=%(ENV_ENABLE_STREAMING)s
startsecs=4
startretries=10
stopwaitsecs=10
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.log

; ── 6. Streaming server (CDP → WebSocket → browser viewer) ────────────────────
[program:streaming]
command=/root/.bun/bin/bun run /app/packages/podman-chrome-streaming/src/index.ts
environment=
    CHROME_HOST="localhost",
    CHROME_DEBUG_PORT="%(ENV_CHROME_DEBUG_PORT)s",
    STREAM_HTTP_PORT="%(ENV_STREAM_HTTP_PORT)s",
    STREAM_WS_PORT="%(ENV_STREAM_WS_PORT)s",
    START_URL="%(ENV_START_URL)s"
priority=50
autorestart=true
autostart=%(ENV_ENABLE_STREAMING)s
startsecs=5
startretries=10
stopwaitsecs=15
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/streaming.log
stderr_logfile=/var/log/supervisor/streaming.log
